# smartgrowth-backend/smartgrowth-backend/deploy/smartgrowth-backend.service
[Unit]
Description=SmartGrowth Django backend (Gunicorn)
After=network.target postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/smartgrowth/smartgrowth-backend/smartgrowth-backend
EnvironmentFile=/var/www/smartgrowth/smartgrowth-backend/smartgrowth-backend/.env
ExecStart=/var/www/smartgrowth/smartgrowth-backend/smartgrowth-backend/.venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8002 \
    config.wsgi:application
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
